# mediamtx.service
# /etc/systemd/system/mediamtx.service
#
# MediaMTX — single static Go binary at /usr/local/bin/mediamtx
# Config: /etc/nexvue/mediamtx.yml

[Unit]
Description=MediaMTX media server (NexVUE)
After=network-online.target
Wants=network-online.target

[Service]
Type=exec
User=nexvue
Group=nexvue
ExecStart=/usr/local/bin/mediamtx /etc/nexvue/mediamtx.yml

Restart=always
RestartSec=2
StartLimitIntervalSec=0

# Allows binding :443 directly in Phase 3 (DMZ) without running as root.
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE

# Hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes
ReadOnlyPaths=/etc/nexvue

[Install]
WantedBy=multi-user.target
